# deploy/private-temp-mail.service
[Unit]
Description=Private Temp Mail Next.js app
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/.openclaw/workspace/private-temp-mail
Environment=NODE_ENV=production
EnvironmentFile=/root/.openclaw/workspace/private-temp-mail/.env.production
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3001
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
